# usr/lib/systemd/system/console-setup.service.d/30_fix.conf
## workaround for:
## console-setup.sh[407]: /bin/setupcon: 866: /bin/setupcon: cannot open /tmp/tmpkbd.5RvhtP: No such file
## https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=846256#44

[Unit]
After=systemd-tmpfiles-setup.service

## TODO: review
## Only run console font/keymap setup when a framebuffer console exists. On a headless or
## serial-only boot (no GPU, hence no /dev/fb0) setupcon's setfont cannot work and spams
## 'setfont: ERROR ioctl(KDFONTOP): Invalid argument' repeatedly. Gating on /dev/fb0 skips the
## service cleanly there (a normal graphical boot has /dev/fb0, so it is unaffected).
#ConditionPathExists=/dev/fb0
