# lib/systemd/system/console-setup.service.d/30_fix.conf
## workaround for:
## console-setup.sh[407]: /bin/setupcon: 866: /bin/setupcon: cannot open /tmp/tmpkbd.5RvhtP: No such file
## https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=846256#44

[Unit]
After=systemd-tmpfiles-setup.service
