# usr/lib/systemd/system/usbguard-dbus.service.d/30_usability-misc.conf
## Copyright (C) 2025 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Service]
ExecCondition=/usr/libexec/usability-misc/has-pci-usb-controller
