# lib/systemd/system/auditd.service.d/30_usability-misc.conf
## This file is part of Whonix.
## Copyright (C) 2016 - 2016 Patrick Schleizer <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
## please use configuration folder /etc/audit/rules.d/ by default
## http://bugs.debian.org/cgi-bin/bugreport.cgi?bug=833474
ExecStartPre=-/sbin/augenrules --load
